# lab3_rules.xml

<group name="windows,sysmon,attack_lab,lolbas">

  <rule id="210100" level="0">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\addinutil\.exe$</field>
    <description>Base AddinUtil.exe via Sysmon on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210101" level="0">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\addinutil\.exe$</field>
    <description>Base AddinUtil.exe via Security Log on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210102" level="12">
    <if_sid>210100,210101</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/](?:addinroot|pipelineroot)\s*:\s*["']?[A-Za-z]:\\(?:users|public|temp|tmp|downloads|desktop|appdata|programdata|windows\\temp|recycler|\$recycle\.bin)(?:\\|["']|$)</field>
    <description>Suspicious AddinUtil.exe: sensitive path argument on $(win.system.computer)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210103" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)\\addinutil\.exe$</field>
    <field name="win.eventdata.image" type="pcre2">(?i)\\(?:bash|bitsadmin|certutil|cmd|cmstp|cscript|installutil|mshta|msiexec|net|powershell|pwsh|regsvr32|rundll32|schtasks|wmic|wscript)\.exe$</field>
    <description>Critical AddinUtil.exe: child process [$(win.eventdata.image)] spawned on $(win.system.computer)</description>
    <mitre><id>T1218</id><id>T1059</id></mitre>
  </rule>

  <rule id="210104" level="15">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.parentProcessName" type="pcre2">(?i)\\addinutil\.exe$</field>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\(?:bash|bitsadmin|certutil|cmd|cmstp|cscript|installutil|mshta|msiexec|net|powershell|pwsh|regsvr32|rundll32|schtasks|wmic|wscript)\.exe$</field>
    <description>Critical AddinUtil.exe: child process [$(win.eventdata.newProcessName)] spawned on $(win.system.computer)</description>
    <mitre><id>T1218</id><id>T1059</id></mitre>
  </rule>

  <rule id="210105" level="13">
    <if_sid>61605</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\addinutil\.exe$</field>
    <field name="win.eventdata.destinationIp" type="pcre2">^(?!10\.|192\.168\.|172\.(?:1[6-9]|2[0-9]|3[01])\.|127\.|169\.254\.|0\.0\.0\.0)</field>
    <description>High AddinUtil.exe: outbound connection to external IP $(win.eventdata.destinationIp):$(win.eventdata.destinationPort)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210106" level="10">
    <if_sid>61613</if_sid>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)\\Addins\.store$</field>
    <description>Suspicious AddinUtil.exe: Addins.store created at [$(win.eventdata.targetFilename)]</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210107" level="13">
    <if_sid>210106</if_sid>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)[A-Za-z]:\\(?:users|public|temp|tmp|downloads|desktop|appdata|programdata|windows\\temp)\\</field>
    <description>High AddinUtil.exe: Addins.store in sensitive directory [$(win.eventdata.targetFilename)]</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210108" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\ysoserial\.exe$</field>
    <description>Critical: ysoserial payload generator detected via Sysmon (Image) on $(win.system.computer)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210109" level="15">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\ysoserial\.exe$</field>
    <description>Critical: ysoserial payload generator detected via Security Log (ProcessName) on $(win.system.computer)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210110" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.originalFileName" type="pcre2">(?i)ysoserial\.exe$</field>
    <description>Critical: ysoserial detected via Sysmon OriginalFileName (renamed binary evasion) on $(win.system.computer)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210111" level="15">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.originalFileName" type="pcre2">(?i)ysoserial\.exe$</field>
    <description>Critical: ysoserial detected via Security Log OriginalFileName (renamed binary evasion) on $(win.system.computer)</description>
    <mitre><id>T1218</id></mitre>
  </rule>

  <rule id="210200" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^4104$</field>
    <field name="win.eventdata.scriptBlockText" type="pcre2">(?i)ms-appinstaller:\/\/(?:\?source=)?https?:</field>
    <description>Medium AppInstaller.exe: ms-appinstaller URI in PowerShell ScriptBlock on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210201" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\appinstaller\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)ms-appinstaller:\/\/(?:\?source=)?https?:</field>
    <description>Suspicious AppInstaller.exe: ms-appinstaller URI execution on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210202" level="12">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\appinstaller\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)ms-appinstaller:\/\/(?:\?source=)?https?:</field>
    <description>Suspicious AppInstaller.exe: ms-appinstaller URI execution on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210203" level="10">
    <if_sid>61613</if_sid>
    <field name="win.eventdata.targetFilename" type="pcre2">(?i)\\Microsoft\.DesktopAppInstaller_[^\\]+\\AC\\INetCache\\</field>
    <description>Suspicious AppInstaller.exe: file cached in INetCache on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210204" level="10">
    <if_sid>61605</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\appinstaller\.exe$</field>
    <field name="win.eventdata.destinationIp" type="pcre2">^(?!10\.|192\.168\.|172\.(?:1[6-9]|2[0-9]|3[01])\.|127\.|169\.254\.|0\.0\.0\.0)</field>
    <description>Medium AppInstaller.exe: external network connection to $(win.eventdata.destinationIp)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210205" level="14">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)\\appinstaller\.exe$</field>
    <field name="win.eventdata.image" type="pcre2">(?i)\\(?:cmd|cscript|mshta|powershell|pwsh|regsvr32|rundll32|wscript|certutil|msiexec)\.exe$</field>
    <description>High AppInstaller.exe: suspicious child process [$(win.eventdata.image)] spawned on $(win.system.computer)</description>
    <mitre><id>T1105</id><id>T1059</id></mitre>
  </rule>

  <rule id="210206" level="14">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.parentProcessName" type="pcre2">(?i)\\appinstaller\.exe$</field>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\(?:cmd|cscript|mshta|powershell|pwsh|regsvr32|rundll32|wscript|certutil|msiexec)\.exe$</field>
    <description>High AppInstaller.exe: suspicious child process [$(win.eventdata.newProcessName)] spawned on $(win.system.computer)</description>
    <mitre><id>T1105</id><id>T1059</id></mitre>
  </rule>

  <rule id="210300" level="5">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\at\.exe$</field>
    <description>At.exe (deprecated scheduler) executed via Sysmon on $(win.system.computer)</description>
    <mitre><id>T1053.002</id></mitre>
    <options>no_full_log</options>
  </rule>

  <rule id="210301" level="5">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\at\.exe$</field>
    <description>At.exe (deprecated scheduler) executed via Security Log on $(win.system.computer)</description>
    <mitre><id>T1053.002</id></mitre>
    <options>no_full_log</options>
  </rule>

  <rule id="210302" level="12">
    <if_sid>210300,210301</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?:[-\/]interactive\b|\\\\[A-Za-z0-9][A-Za-z0-9._-]*)</field>
    <description>Suspicious At.exe: interactive or remote scheduled task on $(win.system.computer)</description>
    <mitre><id>T1053.002</id></mitre>
  </rule>

  <rule id="210303" level="10">
    <if_sid>210300,210301</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)\b\d{1,2}:\d{2}(?:am|pm)?\b</field>
    <description>Suspicious At.exe: new task scheduled with time argument on $(win.system.computer)</description>
    <mitre><id>T1053.002</id></mitre>
  </rule>

  <rule id="210400" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\bitsadmin\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/](?:transfer|addfile)\b.+(?:https?|ftp|smb):\/\/</field>
    <description>Suspicious Bitsadmin.exe: remote file download via BITS on $(win.system.computer)</description>
    <mitre><id>T1197</id><id>T1105</id></mitre>
  </rule>

  <rule id="210401" level="12">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\bitsadmin\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/](?:transfer|addfile)\b.+(?:https?|ftp|smb):\/\/</field>
    <description>Suspicious Bitsadmin.exe: remote file download via BITS on $(win.system.computer)</description>
    <mitre><id>T1197</id><id>T1105</id></mitre>
  </rule>

  <rule id="210402" level="14">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\bitsadmin\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/]setnotifycmdline\b</field>
    <description>High Bitsadmin.exe: persistence via SetNotifyCmdLine on $(win.system.computer)</description>
    <mitre><id>T1197</id><id>T1546</id></mitre>
  </rule>

  <rule id="210403" level="14">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\bitsadmin\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/]setnotifycmdline\b</field>
    <description>High Bitsadmin.exe: persistence via SetNotifyCmdLine on $(win.system.computer)</description>
    <mitre><id>T1197</id><id>T1546</id></mitre>
  </rule>

  <rule id="210500" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\certreq\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?=.*[-\/]post\b)(?=.*[-\/]config\b)</field>
    <description>Suspicious CertReq.exe: -Post -Config download on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210501" level="12">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\certreq\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?=.*[-\/]post\b)(?=.*[-\/]config\b)</field>
    <description>Suspicious CertReq.exe: -Post -Config download on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210502" level="13">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\certreq\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?=.*[-\/]post\b)(?=.*https?:\/\/)</field>
    <description>Suspicious CertReq.exe: HTTP POST to remote URL - potential exfiltration on $(win.system.computer)</description>
    <mitre><id>T1105</id><id>T1048</id></mitre>
  </rule>

  <rule id="210503" level="13">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\certreq\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?=.*[-\/]post\b)(?=.*https?:\/\/)</field>
    <description>Suspicious CertReq.exe: HTTP POST to remote URL - potential exfiltration on $(win.system.computer)</description>
    <mitre><id>T1105</id><id>T1048</id></mitre>
  </rule>

  <rule id="210600" level="0">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\certutil\.exe$</field>
    <description>Base CertUtil.exe via Sysmon on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210601" level="0">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\certutil\.exe$</field>
    <description>Base CertUtil.exe via Security Log on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210602" level="13">
    <if_sid>210600,210601</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/]urlcache\b</field>
    <description>High CertUtil.exe: -urlcache remote file download on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210602b" level="10">
    <if_sid>210600,210601</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/](?:verifyurl|verifyctl|url)\b</field>
    <description>Medium CertUtil.exe: remote URL verification/download on $(win.system.computer)</description>
    <mitre><id>T1105</id></mitre>
  </rule>

  <rule id="210603" level="13">
    <if_sid>210600,210601</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/](?:encode|encodehex|decode|decodehex)\b</field>
    <description>High CertUtil.exe: encode/decode operation (payload staging) on $(win.system.computer)</description>
    <mitre><id>T1140</id></mitre>
  </rule>

  <rule id="210700" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\ftp\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/]s:</field>
    <description>Suspicious Ftp.exe: script file execution on $(win.system.computer)</description>
    <mitre><id>T1202</id><id>T1105</id></mitre>
  </rule>

  <rule id="210701" level="12">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\ftp\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[-\/]s:</field>
    <description>Suspicious Ftp.exe: script file execution on $(win.system.computer)</description>
    <mitre><id>T1202</id><id>T1105</id></mitre>
  </rule>

  <rule id="210702" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)\\ftp\.exe$</field>
    <field name="win.eventdata.image" type="pcre2">(?i)\\(?:cmd|powershell|pwsh)\.exe$</field>
    <description>Critical Ftp.exe: shell [$(win.eventdata.image)] spawned on $(win.system.computer)</description>
    <mitre><id>T1202</id><id>T1059</id></mitre>
  </rule>

  <rule id="210703" level="15">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.parentProcessName" type="pcre2">(?i)\\ftp\.exe$</field>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\(?:cmd|powershell|pwsh)\.exe$</field>
    <description>Critical Ftp.exe: shell spawned on $(win.system.computer)</description>
    <mitre><id>T1202</id><id>T1059</id></mitre>
  </rule>

  <rule id="210800" level="0">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)\\hh\.exe$</field>
    <description>Base Hh.exe via Sysmon on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210801" level="0">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\hh\.exe$</field>
    <description>Base Hh.exe via Security Log on $(win.system.computer)</description>
    <options>no_full_log</options>
  </rule>

  <rule id="210802" level="12">
    <if_sid>210800,210801</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(?=.*\.chm)(?=.*(?:https?:|\\\\[A-Za-z0-9][A-Za-z0-9._-]*))</field>
    <description>Suspicious Hh.exe: remote CHM via HTTP or UNC path on $(win.system.computer)</description>
    <mitre><id>T1218.001</id><id>T1105</id></mitre>
  </rule>

  <rule id="210802b" level="11">
    <if_sid>210800,210801</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)[A-Za-z]:\\(?:users|public|temp|tmp|downloads|appdata|programdata|windows\\temp)\\[^\\]+\.chm</field>
    <description>Suspicious Hh.exe: CHM file loaded from suspicious local path on $(win.system.computer)</description>
    <mitre><id>T1218.001</id></mitre>
  </rule>

  <rule id="210803" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.parentImage" type="pcre2">(?i)\\hh\.exe$</field>
    <field name="win.eventdata.image" type="pcre2">(?i)\\(?:cmd|cscript|mshta|powershell|pwsh|regsvr32|rundll32|wscript)\.exe$</field>
    <description>Critical Hh.exe: child process [$(win.eventdata.image)] spawned on $(win.system.computer)</description>
    <mitre><id>T1218.001</id><id>T1059</id></mitre>
  </rule>

  <rule id="210804" level="15">
    <if_sid>60109</if_sid>
    <field name="win.eventdata.parentProcessName" type="pcre2">(?i)\\hh\.exe$</field>
    <field name="win.eventdata.newProcessName" type="pcre2">(?i)\\(?:cmd|cscript|mshta|powershell|pwsh|regsvr32|rundll32|wscript)\.exe$</field>
    <description>Critical Hh.exe: child process spawned on $(win.system.computer)</description>
    <mitre><id>T1218.001</id><id>T1059</id></mitre>
  </rule>

</group>